# common/systemd/katafract-update.service
[Unit]
Description=Katafract node agent self-update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/katafract-update
StandardOutput=journal
StandardError=journal
